# uwsgi_start.ini
[program:uwsgi_start]

directory = /home/gz_v1.0.0
command = /root/.virtualenvs/django_gz/bin/uwsgi --ini /home/gz_v1.0.0/uwsgi.ini

user = root
numprocs = 1
stdout_logfile = /home/gz_v1.0.0/log/uwsgi_start.log
stderr_logfile = /home/gz_v1.0.0/log/uwsgi_start.log
autostart = true
autorestart = true
startsecs = 50
