# uwsgi_start.ini
[program:uwsgi_start]

directory = /home/gz_v1.0.0
command = /root/.virtualenvs/gz/bin/uwsgi --ini /home/gz_v1.0.0/uwsgi.ini

user = root


stdout_logfile = /home/gz_v1.0.0/log/uwsgi_out.log
stderr_logfile = /home/gz_v1.0.0/log/uwsgi_err.log
autostart = true
autorestart = true
startsecs = 0
stopwaitsecs = 0
stopsignal = QUIT
redirect_stderr = true
loglevel = error
logfile_maxbytes = 1M
